# hyprland/.config/systemd/user/theme-light.service
[Unit]
Description=Switch to light theme

[Service]
Type=oneshot
ExecStart=%h/.local/bin/theme-light
